# crontab
30 1 * * * {{MFSERV_HOME}}/bin/cronwrap.sh -l -e -- plugin_wrapper force_common_files force.sh >>{{MODULE_RUNTIME_HOME}}/log/force.log 2>&1
